# scripts/com.kumax.openclaw-backup.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.kumax.openclaw-backup.daily</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/kumax/.openclaw/workspace/scripts/backup-openclaw-daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>3</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/kumax/.openclaw/workspace/tmp/openclaw-backup-launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kumax/.openclaw/workspace/tmp/openclaw-backup-launchd.err.log</string>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
